# go2w_remote_arm/systemd/go2w-web-monitor.service
[Unit]
Description=Go2W Web Monitor — LAN dashboard (status/nodes/topics/camera/teleop)
After=network-online.target multi-user.target
Wants=network-online.target

[Service]
Type=simple
User=unitree
Group=dialout
# HOME so ros2 finds /home/unitree/.ros logs; dialout not strictly needed here
# but kept for parity with the launcher (harmless).
Environment=HOME=/home/unitree
# /bin/bash -lc sources /etc/profile.d/*.sh so any CycloneDDS / FastDDS env
# the user configured for the Go2W bridge is picked up (same domain as teleop).
ExecStart=/bin/bash -lc 'source /opt/ros/humble/setup.bash && source /home/unitree/ros2_ws/install/setup.bash && exec ros2 run go2w_remote_arm web_monitor'
KillMode=control-group
KillSignal=SIGINT
TimeoutStopSec=10
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
