# om6dof/om6dof_teleop/systemd/om6dof-web-monitor.service
[Unit]
Description=Go2W Web Monitor — LAN dashboard (status/nodes/topics/camera/teleop)
After=network-online.target multi-user.target
Wants=network-online.target

[Service]
Type=simple
User=unitree
Group=dialout
# HOME so ros2 finds /home/unitree/.ros logs; dialout not strictly needed here
# but kept for parity with the launcher (harmless).
Environment=HOME=/home/unitree
# Use the same Unitree CycloneDDS environment and ROS domain as the hardware
# stack, then overlay the local workspace.
ExecStart=/bin/bash -c 'set -e; source /home/unitree/unitree_ros2/setup.sh; source /home/unitree/ros2_ws/install/setup.bash; exec ros2 run om6dof_teleop web_monitor'
KillMode=control-group
KillSignal=SIGINT
TimeoutStopSec=10
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
